# quick-core/src/main/resources/doc/supervistor.ini
[program:xiuxian]
directory=/opt/services/spring-boot-quick
command=sh /opt/services/spring-boot-quick/develop.sh 8080
autostart=true
autorestart=true
priority=5
user=root
stdout_logfile=/dev/null
redirect_stderr=true
loglevel=info
stopasgroup=true
killasgroup=true

[program:xiuxian1]
directory=/opt/services/spring-boot-quick
command=sh /opt/services/spring-boot-quick/develop.sh 8081
autostart=true
autorestart=true
priority=5
user=root
stdout_logfile=/dev/null
redirect_stderr=true
loglevel=info
stopasgroup=true
killasgroup=true
